# deploy/supervisor/ourth-worker.conf
; /etc/supervisor/conf.d/ourth-worker.conf
;
; Supervisor configuration for the Ourth Laravel queue worker.
; Install: sudo supervisorctl reread && sudo supervisorctl update
; Start:   sudo supervisorctl start ourth-worker:*

[program:ourth-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/ourth-app/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/ourth-worker.log
stopwaitsecs=3600
